# 99-monitor-hotplug.rules.in
# udev rule for automatic display positioning
# Place this file in /etc/udev/rules.d/
#
# This rule triggers when a display is connected via DisplayPort, HDMI, or VGA

# Trigger on DRM (Direct Rendering Manager) change events
ACTION=="change", SUBSYSTEM=="drm", RUN+="@SCRIPT_DIR@/x11-auto.udev.sh '%k' '$env{DEVNAME}' '$env{DEVPATH}'"


